# recipes-tests/holoscan-test-suite/files/services/holoscan-test-suite-clara-agx-devkit.service
[Unit]
Description="Holoscan Test Suite: Clara AGX DevKit"

[Service]
Type=simple
WorkingDirectory=/opt/nvidia/holoscan-test-suite
ExecStart=/usr/bin/screen -D -m -S holoscan-test-suite /bin/sh run_clara_agx_devkit.sh
ExecStop=/usr/bin/screen -X -S holoscan-test-suite
Restart=always
RestartSec=30sec

[Install]
WantedBy=multi-user.target
